# rootfs/etc/supervisor/1-singbox.conf
[program:singbox]
command=bash /opt/singbox.sh
priority=1
autorestart=unexpected
stderr_logfile=/proc/self/fd/2
