# rootfs/etc/supervisor/1-xvfb.conf
[program:Xvfb]
command=bash -c 'Xvfb $DISPLAY -screen 0 $VNC_RESOLUTION"x24" -nolisten unix'
priority=1
stderr_logfile=/proc/self/fd/2
