# rootfs/etc/supervisor/3-vnc.conf
[program:VNC]
command=bash /opt/x11vnc.sh
priority=3
startsecs=0
autorestart=unexpected
stderr_logfile=/proc/self/fd/2
