# rootfs/etc/supervisor/5-grass.conf
[program:Grass]
command=grass
priority=5
autorestart=unexpected
stderr_logfile=/proc/self/fd/2
