# rootfs/etc/supervisor/6-mihomo.conf
[program:mihomo]
command=bash /opt/mihomo.sh
priority=6
autorestart=unexpected
stderr_logfile=/proc/self/fd/2
