# rootfs/etc/supervisor/7-grass.conf
[program:Grass]
command=grass
priority=7
autorestart=unexpected
stderr_logfile=/proc/self/fd/2
